# build/etc/cron.d/packageiq
# run piq_rhel every day at 12:30AM

30 0 * * * root /opt/packageiq/bin/piq_rhel