# build/etc/profile.d/packageiq.sh
# add the packageiq omnibus bin dir to pathmunge

pathmunge /opt/packageiq/bin